# deploy/vps/eurbanizam-api.service
[Unit]
Description=eUrbanizam API
After=network.target

[Service]
Type=simple
User=niki
WorkingDirectory=/home/niki/eurbanizam-app
Environment=API_HOST=127.0.0.1
Environment=API_PORT=8000
Environment=API_RELOAD=false
Environment=CORS_ALLOW_ORIGINS=http://127.0.0.1:3000,http://5.189.136.118:3000,http://5.189.136.118,https://eurbanizam.easy.mk
EnvironmentFile=-/home/niki/.eurbanizam/secrets/.eurbanizam_secrets.env
ExecStart=/home/niki/eurbanizam-app/.venv/bin/python -m uvicorn apps.api.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
