# vps-landing/eurbanizam-healthcheck.timer
[Unit]
Description=Run Eurbanizam healthcheck every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true
Unit=eurbanizam-healthcheck.service

[Install]
WantedBy=timers.target
