# vps-landing/eurbanizam-daily.timer
[Unit]
Description=Run eUrbanizam daily sync and report at 04:00

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true
Unit=eurbanizam-daily.service

[Install]
WantedBy=timers.target
